# plugins/protocol_lws_dht_dnssec_monitor/usr-lib-systemd-system-lwsws-dnssec-monitor.service
[Unit]
Description=Libwebsockets Web Server
After=syslog.target

[Service]
ExecStart=/usr/local/bin/lwsws
ExecReload=/usr/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
