# minimal-examples-lowlevel/http-server/minimal-http-server-systemd-socketact/lws-minimal-http-server-systemd-socketact.socket
[Socket]
ListenStream=7681

[Install]
WantedBy=sockets.target

