# minimal-examples-lowlevel/http-server/minimal-http-server-systemd-socketact/lws-minimal-http-server-systemd-socketact.service
[Unit]
Description=LWS minimal example for socket activation

[Service]
ExecStart=/usr/local/bin/lws-minimal-http-server-systemd-socketact
